# scripts/gunicorn-treinacnh.service
[Unit]
Description=Gunicorn daemon for TREINACNH
After=network.target
After=mysql.service

[Service]
User=root
Group=www-data
WorkingDirectory=/var/www/TREINACNH
Environment="PATH=/var/www/TREINACNH/venv/bin"
ExecStart=/var/www/TREINACNH/venv/bin/gunicorn \
          --workers 3 \
          --bind 127.0.0.1:8001 \
          --timeout 120 \
          --access-logfile /var/www/TREINACNH/logs/gunicorn-access.log \
          --error-logfile /var/www/TREINACNH/logs/gunicorn-error.log \
          config.wsgi:application
Restart=always
RestartSec=10
StartLimitIntervalSec=0

[Install]
WantedBy=multi-user.target
